# .config/systemd/user/labwc-session.target
[Unit]
Description=Labwc Session
After=default.target
