# .config/systemd/user/voxtype.service.d/override.conf
[Unit]
PartOf=labwc-session.target
After=labwc-session.target

[Install]
WantedBy=labwc-session.target
